# Library/profile.d/env.sh
export PAGER=less
export EDITOR=vim
export RSYNC_RSH=ssh

export PYTHONSTARTUP=$HOME/Library/pythonrc
export PYTHONPATH=$PYTHONPATH:$HOME/opt/python3.8
export PYTHONPATH=$PYTHONPATH:$HOME/dev/var/libs/pyrtdb

export LD_LIBRARY_PATH=$LD_LIBRARY_PATH:$HOME/opt/lib
export PATH=$PATH:$HOME/opt/bin:$HOME/.local/bin
export PKG_CONFIG_PATH=/usr/share/pkgconfig:$HOME/opt/lib/pkgconfig

#export QT_QPA_PLATFORMTHEME="qt5ct"
#export QT_AUTO_SCREEN_SCALE_FACTOR=0
#export QT_STYLE_OVERRIDE=adwaita-dark

